# deploy/maps-backend.service
[Unit]
Description=Maps Backend API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/maps-backend
ExecStart=/opt/maps-backend/maps-backend
Restart=always
RestartSec=3
EnvironmentFile=/opt/maps-backend/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=maps-backend

[Install]
WantedBy=multi-user.target
